# net-p2p/bittorrent-sync/files/init.d/btsync
#!/sbin/runscript

# Copyright 2013-2014 Jonathan Vasquez <jvasquez1011@gmail.com>
# Distributed under the terms of the GNU General Public License v2

NAME="BitTorrent Sync"
SYNC_NAME=${RC_SVCNAME}
SYNC_BINARY="/opt/btsync/btsync"
SYNC_OPTS="${SYNC_OPTS} --nodaemon --config /etc/btsync/${SYNC_NAME}.conf"
SYNC_PIDFILE="/var/run/btsync/${SYNC_NAME}.pid"
[ -z "${SYNC_UMASK}" ] && SYNC_UMASK="022"
[ -z "${SYNC_USER}" ] && SYNC_USER="btsync"


start_pre() {
	# Ensure the PID file directory exists and is writable
	checkpath --owner ${SYNC_USER}:btsync --mode 0775 --directory /var/run/btsync/
}

start() {
	ebegin "Starting ${NAME}"

	start-stop-daemon --start --exec "${SYNC_BINARY}" \
	   --pidfile "${SYNC_PIDFILE}" --background \
	   --umask "${SYNC_UMASK}" --user "${SYNC_USER}" \
	   -- ${SYNC_OPTS}

	eend $?
}

stop() {
	ebegin "Stopping ${NAME}"

	start-stop-daemon --stop --exec "${SYNC_BINARY}" \
	    --pidfile "${SYNC_PIDFILE}"

	eend $?
}

reload() {
	ebegin "Reloading ${NAME}"

	start-stop-daemon --signal HUP --exec "${SYNC_BINARY}" \
	    --pidfile "${SYNC_PIDFILE}"

	eend $?
}
